# systemd/.config/systemd/user/picom.service
[Unit]
Description=Picom

[Service]
Type=simple
ExecStart=/usr/bin/picom
Restart=on-failure

[Install]
WantedBy=default.target
