# systemd/.config/systemd/user/xscreensaver.service
[Unit]
Description=XScreenSaver

[Service]
Type=simple
ExecStart=xscreensaver -nosplash
Restart=on-failure

[Install]
WantedBy=default.target
